# Library/LaunchAgents/dsaenztagarro.github.io.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dsaenztagarro.github.io</string>
    <key>RunAtLoad</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Disabled</key>
    <false/>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/dst/Scripts/dsaenztagarro.github.io.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/dst/Library/Logs/dsaenztagarro.github.io.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/dst/Library/Logs/dsaenztagarro.github.io.err</string>
    <key>Debug</key>
    <true/>
</dict>
</plist>